# deploy/monita-collector.service
[Unit]
Description=Monita Collector API
After=network.target

[Service]
Type=simple
User=monita
Group=monita
WorkingDirectory=/etc/monita-collector
EnvironmentFile=/etc/monita-collector/.env
ExecStart=/usr/local/bin/monita-collector
Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/etc/monita-collector

[Install]
WantedBy=multi-user.target
